# tools/launchd/com.phototestapp.portfolio-prices.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.phototestapp.portfolio-prices</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "__REPO_ROOT__" &amp;&amp; FETCH_MUTUAL_FUNDS=1 FETCH_JAPAN_STOCKS=1 FETCH_CRYPTO=1 AUTO_PUBLISH=1 bash tools/publish_prices.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>5</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/portfolio-prices.out.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/portfolio-prices.err.log</string>
</dict>
</plist>
